# pm2-mindsetai.service
[Unit]
Description=PM2 process manager for mindsetai-api
Documentation=https://pm2.keymetrics.io/
After=network.target

[Service]
Type=simple
User=%i
Environment=PATH=/opt/alt/alt-nodejs20/root/usr/bin:/usr/local/bin:/usr/bin:/bin
Environment=PM2_HOME=/home/%i/.pm2
ExecStart=/usr/local/bin/pm2 resurrect
ExecReload=/usr/local/bin/pm2 reload all
ExecStop=/usr/local/bin/pm2 kill
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
